# modules/disconnect-lock/10-disconnect-lock.rules
ATTRS{idVendor}=="13fe", ATTRS{idProduct}=="1d00" ACTION=="remove" RUN+="/usr/bin/disconnect-lock"
